# systemd/autotune.service
[Unit]
Description=Autotune Self-Improvement Pipeline
After=network-online.target

[Service]
Type=oneshot
Environment=NODE_ENV=production
Environment=AUTOTUNE_CONFIG_PATH=%h/.config/dotclaw/autotune.json
EnvironmentFile=-%h/.config/dotclaw/autotune.env
ExecStart=/usr/bin/env autotune once

[Install]
WantedBy=multi-user.target
